# tools/trade-alerts/com.cortana.canslim-precompute.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortana.canslim-precompute</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/hd/openclaw/tools/trade-alerts/canslim-precompute.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/hd/openclaw/tools/trade-alerts</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>

    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>12</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>

    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>5</integer>
      <key>Hour</key>
      <integer>15</integer>
      <key>Minute</key>
      <integer>25</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/hd/.openclaw/logs/canslim-precompute.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/hd/.openclaw/logs/canslim-precompute.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
