# config/launchd/com.cortana.oracle-precompute.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.cortana.oracle-precompute</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>python3</string>
      <string>/Users/hd/openclaw/tools/oracle/precompute.py</string>
      <string>run</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>5</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/hd/openclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/opt/postgresql@17/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/hd/openclaw/tmp/oracle-precompute.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hd/openclaw/tmp/oracle-precompute.stderr.log</string>
  </dict>
</plist>
